# infra/service-install/systemd.service.tmpl
[Unit]
Description=Bematist telemetry collector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart={{BINARY_PATH}} run
Restart=on-failure
RestartSec=5s
StandardOutput=append:{{HOME}}/.bm-pilot/bm-pilot.log
StandardError=append:{{HOME}}/.bm-pilot/bm-pilot.log

[Install]
WantedBy=default.target
